# recipes-scripts/ntpd-service/ntpd-service-1.0/ntpd
### BEGIN INIT INFO
# Provides:          ntpd
# Required-Start:    $local_fs $network
# Required-Stop:     $local_fs
# Default-Start:     3 4 5
# Default-Stop:      0 1 2 6
# Short-Description: Set system clock from ntpd
# Description:       Set system clock to network time.
### END INIT INFO

NTP=`which ntpd`
[ -x $NTP ] || NTP=/usr/sbin/ntpd
[ -x $NTP ] || NTP=/usr/local/bin/ntpd
[ -x $NTP ] || exit 0

case "$1" in
        start)
                echo Starting NTP...
                mv /var/log/ntp.log /var/log/ntp.old || true
                tail -n9 /var/log/ntp-start.log > /var/log/ntp-start.log || true
                $( $NTP -g -q -x >> /var/log/ntp-start.log; $NTP -p /var/run/ntpd.pid -I eth0 ) &
                ;;
        stop)
                if [ -f /var/run/ntpd.pid ]
                then    kill `cat /var/run/ntpd.pid`
                        rm /var/run/ntpd.pid
                fi
                ;;
        show)
                ;;
        *)
                echo "Usage: ntpd.sh {start|stop|show|reload|restart}" >&2
                echo "       start sets kernel (system) clock from network (NTP) clock" >&2
                echo "       stop stops ntpd" >&2
                exit 1
                ;;
esac
