# recipes-devtools/python/files/protocolbridge.service
[Unit]
Description=protocolbridge

[Service]
ExecStart=
Type=simple
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
