# recipes-configuration/arista-usb-automount/files/90-usb-storage.rules
KERNEL=="sd[b-z]*", ENV{ID_FS_TYPE}=="vfat", RUN+="/usr/bin/usbmount"